# imooc/conf/supervisor_celery_flower.ini
[program:imooc-celery-flower]
command=python manage.py celery flower --basic_auth=ten:ten
directory=/root/celery_with_django/imooc
environment=PATH="/root/.local/share/virtualenvs/celery_with_django-oEbuwveD/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/root/bin"
stdout_logfile=/root/celery_with_django/imooc/logs/celery.flower.log
stderr_logfile=/root/celery_with_django/imooc/logs/celery.flower.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=996




