# imooc/conf/supervisor_celery_beat.ini
[program:imooc-celery-beat]
command=python manage.py celery beat -l INFO
directory=/root/celery_with_django/imooc
environment=PATH="/root/.local/share/virtualenvs/celery_with_django-oEbuwveD/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/root/bin"
stdout_logfile=/root/celery_with_django/imooc/logs/celery.beat.log
stderr_logfile=/root/celery_with_django/imooc/logs/celery.beat.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=997




