# imooc/conf/supervisor_celery_worker.ini
[program:imooc-celery-worker]
command=python manage.py celery worker -l INFO
directory=/root/celery_with_django/imooc
environment=PATH="/root/.local/share/virtualenvs/celery_with_django-oEbuwveD/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/root/bin"
stdout_logfile=/root/celery_with_django/imooc/logs/celery.worker.log
stderr_logfile=/root/celery_with_django/imooc/logs/celery.worker.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=998




